# configs/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock; use a unix:// URL  for a unix socket

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:apache2]
autorestart=true
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:mysqld]
autorestart=true
command=/bin/bash -c "/usr/bin/mysqld_safe"
